# conf/agora-pusher.service
[Unit]
Description=Agora Garden Pusher (Syncs hosted gardens to Forgejo)
After=network.target

[Service]
WorkingDirectory=%h/agora-bridge
ExecStart=%h/agora-bridge/push-gardens.sh
Restart=always
RestartSec=60

[Install]
WantedBy=default.target
